# reboot-schedule.service
[Unit]
Description=Scheduled reboot service

[Service]
Type=oneshot
ExecStart=/sbin/reboot
