# timelapse.timer
[Unit]
Description=Run timelapse script, every 2 min during 4:00-8:00, 16:00-20:30; and 12:56, 13:03

[Timer]
# 04:00–08:00 每2分钟执行一次
OnCalendar=*-*-* 04:00
OnCalendar=*-*-* 04:02
OnCalendar=*-*-* 04:04
OnCalendar=*-*-* 04:06
OnCalendar=*-*-* 04:08
OnCalendar=*-*-* 04:10
OnCalendar=*-*-* 04:12
OnCalendar=*-*-* 04:14
OnCalendar=*-*-* 04:16
OnCalendar=*-*-* 04:18
OnCalendar=*-*-* 04:20
OnCalendar=*-*-* 04:22
OnCalendar=*-*-* 04:24
OnCalendar=*-*-* 04:26
OnCalendar=*-*-* 04:28
OnCalendar=*-*-* 04:30
OnCalendar=*-*-* 04:32
OnCalendar=*-*-* 04:34
OnCalendar=*-*-* 04:36
OnCalendar=*-*-* 04:38
OnCalendar=*-*-* 04:40
OnCalendar=*-*-* 04:42
OnCalendar=*-*-* 04:44
OnCalendar=*-*-* 04:46
OnCalendar=*-*-* 04:48
OnCalendar=*-*-* 04:50
OnCalendar=*-*-* 04:52
OnCalendar=*-*-* 04:54
OnCalendar=*-*-* 04:56
OnCalendar=*-*-* 04:58
OnCalendar=*-*-* 05:00
OnCalendar=*-*-* 05:02
OnCalendar=*-*-* 05:04
OnCalendar=*-*-* 05:06
OnCalendar=*-*-* 05:08
OnCalendar=*-*-* 05:10
OnCalendar=*-*-* 05:12
OnCalendar=*-*-* 05:14
OnCalendar=*-*-* 05:16
OnCalendar=*-*-* 05:18
OnCalendar=*-*-* 05:20
OnCalendar=*-*-* 05:22
OnCalendar=*-*-* 05:24
OnCalendar=*-*-* 05:26
OnCalendar=*-*-* 05:28
OnCalendar=*-*-* 05:30
OnCalendar=*-*-* 05:32
OnCalendar=*-*-* 05:34
OnCalendar=*-*-* 05:36
OnCalendar=*-*-* 05:38
OnCalendar=*-*-* 05:40
OnCalendar=*-*-* 05:42
OnCalendar=*-*-* 05:44
OnCalendar=*-*-* 05:46
OnCalendar=*-*-* 05:48
OnCalendar=*-*-* 05:50
OnCalendar=*-*-* 05:52
OnCalendar=*-*-* 05:54
OnCalendar=*-*-* 05:56
OnCalendar=*-*-* 05:58
OnCalendar=*-*-* 06:00
OnCalendar=*-*-* 06:02
OnCalendar=*-*-* 06:04
OnCalendar=*-*-* 06:06
OnCalendar=*-*-* 06:08
OnCalendar=*-*-* 06:10
OnCalendar=*-*-* 06:12
OnCalendar=*-*-* 06:14
OnCalendar=*-*-* 06:16
OnCalendar=*-*-* 06:18
OnCalendar=*-*-* 06:20
OnCalendar=*-*-* 06:22
OnCalendar=*-*-* 06:24
OnCalendar=*-*-* 06:26
OnCalendar=*-*-* 06:28
OnCalendar=*-*-* 06:30
OnCalendar=*-*-* 06:32
OnCalendar=*-*-* 06:34
OnCalendar=*-*-* 06:36
OnCalendar=*-*-* 06:38
OnCalendar=*-*-* 06:40
OnCalendar=*-*-* 06:42
OnCalendar=*-*-* 06:44
OnCalendar=*-*-* 06:46
OnCalendar=*-*-* 06:48
OnCalendar=*-*-* 06:50
OnCalendar=*-*-* 06:52
OnCalendar=*-*-* 06:54
OnCalendar=*-*-* 06:56
OnCalendar=*-*-* 06:58
OnCalendar=*-*-* 07:00
OnCalendar=*-*-* 07:02
OnCalendar=*-*-* 07:04
OnCalendar=*-*-* 07:06
OnCalendar=*-*-* 07:08
OnCalendar=*-*-* 07:10
OnCalendar=*-*-* 07:12
OnCalendar=*-*-* 07:14
OnCalendar=*-*-* 07:16
OnCalendar=*-*-* 07:18
OnCalendar=*-*-* 07:20
OnCalendar=*-*-* 07:22
OnCalendar=*-*-* 07:24
OnCalendar=*-*-* 07:26
OnCalendar=*-*-* 07:28
OnCalendar=*-*-* 07:30
OnCalendar=*-*-* 07:32
OnCalendar=*-*-* 07:34
OnCalendar=*-*-* 07:36
OnCalendar=*-*-* 07:38
OnCalendar=*-*-* 07:40
OnCalendar=*-*-* 07:42
OnCalendar=*-*-* 07:44
OnCalendar=*-*-* 07:46
OnCalendar=*-*-* 07:48
OnCalendar=*-*-* 07:50
OnCalendar=*-*-* 07:52
OnCalendar=*-*-* 07:54
OnCalendar=*-*-* 07:56
OnCalendar=*-*-* 07:58
OnCalendar=*-*-* 08:00

# 16:00–20:30 每2分钟执行一次
OnCalendar=*-*-* 16:00
OnCalendar=*-*-* 16:02
OnCalendar=*-*-* 16:04
OnCalendar=*-*-* 16:06
OnCalendar=*-*-* 16:08
OnCalendar=*-*-* 16:10
OnCalendar=*-*-* 16:12
OnCalendar=*-*-* 16:14
OnCalendar=*-*-* 16:16
OnCalendar=*-*-* 16:18
OnCalendar=*-*-* 16:20
OnCalendar=*-*-* 16:22
OnCalendar=*-*-* 16:24
OnCalendar=*-*-* 16:26
OnCalendar=*-*-* 16:28
OnCalendar=*-*-* 16:30
OnCalendar=*-*-* 16:32
OnCalendar=*-*-* 16:34
OnCalendar=*-*-* 16:36
OnCalendar=*-*-* 16:38
OnCalendar=*-*-* 16:40
OnCalendar=*-*-* 16:42
OnCalendar=*-*-* 16:44
OnCalendar=*-*-* 16:46
OnCalendar=*-*-* 16:48
OnCalendar=*-*-* 16:50
OnCalendar=*-*-* 16:52
OnCalendar=*-*-* 16:54
OnCalendar=*-*-* 16:56
OnCalendar=*-*-* 16:58
OnCalendar=*-*-* 17:00
OnCalendar=*-*-* 17:02
OnCalendar=*-*-* 17:04
OnCalendar=*-*-* 17:06
OnCalendar=*-*-* 17:08
OnCalendar=*-*-* 17:10
OnCalendar=*-*-* 17:12
OnCalendar=*-*-* 17:14
OnCalendar=*-*-* 17:16
OnCalendar=*-*-* 17:18
OnCalendar=*-*-* 17:20
OnCalendar=*-*-* 17:22
OnCalendar=*-*-* 17:24
OnCalendar=*-*-* 17:26
OnCalendar=*-*-* 17:28
OnCalendar=*-*-* 17:30
OnCalendar=*-*-* 17:32
OnCalendar=*-*-* 17:34
OnCalendar=*-*-* 17:36
OnCalendar=*-*-* 17:38
OnCalendar=*-*-* 17:40
OnCalendar=*-*-* 17:42
OnCalendar=*-*-* 17:44
OnCalendar=*-*-* 17:46
OnCalendar=*-*-* 17:48
OnCalendar=*-*-* 17:50
OnCalendar=*-*-* 17:52
OnCalendar=*-*-* 17:54
OnCalendar=*-*-* 17:56
OnCalendar=*-*-* 17:58
OnCalendar=*-*-* 18:00
OnCalendar=*-*-* 18:02
OnCalendar=*-*-* 18:04
OnCalendar=*-*-* 18:06
OnCalendar=*-*-* 18:08
OnCalendar=*-*-* 18:10
OnCalendar=*-*-* 18:12
OnCalendar=*-*-* 18:14
OnCalendar=*-*-* 18:16
OnCalendar=*-*-* 18:18
OnCalendar=*-*-* 18:20
OnCalendar=*-*-* 18:22
OnCalendar=*-*-* 18:24
OnCalendar=*-*-* 18:26
OnCalendar=*-*-* 18:28
OnCalendar=*-*-* 18:30
OnCalendar=*-*-* 18:32
OnCalendar=*-*-* 18:34
OnCalendar=*-*-* 18:36
OnCalendar=*-*-* 18:38
OnCalendar=*-*-* 18:40
OnCalendar=*-*-* 18:42
OnCalendar=*-*-* 18:44
OnCalendar=*-*-* 18:46
OnCalendar=*-*-* 18:48
OnCalendar=*-*-* 18:50
OnCalendar=*-*-* 18:52
OnCalendar=*-*-* 18:54
OnCalendar=*-*-* 18:56
OnCalendar=*-*-* 18:58
OnCalendar=*-*-* 19:00
OnCalendar=*-*-* 19:02
OnCalendar=*-*-* 19:04
OnCalendar=*-*-* 19:06
OnCalendar=*-*-* 19:08
OnCalendar=*-*-* 19:10
OnCalendar=*-*-* 19:12
OnCalendar=*-*-* 19:14
OnCalendar=*-*-* 19:16
OnCalendar=*-*-* 19:18
OnCalendar=*-*-* 19:20
OnCalendar=*-*-* 19:22
OnCalendar=*-*-* 19:24
OnCalendar=*-*-* 19:26
OnCalendar=*-*-* 19:28
OnCalendar=*-*-* 19:30
OnCalendar=*-*-* 19:32
OnCalendar=*-*-* 19:34
OnCalendar=*-*-* 19:36
OnCalendar=*-*-* 19:38
OnCalendar=*-*-* 19:40
OnCalendar=*-*-* 19:42
OnCalendar=*-*-* 19:44
OnCalendar=*-*-* 19:46
OnCalendar=*-*-* 19:48
OnCalendar=*-*-* 19:50
OnCalendar=*-*-* 19:52
OnCalendar=*-*-* 19:54
OnCalendar=*-*-* 19:56
OnCalendar=*-*-* 19:58
OnCalendar=*-*-* 20:00
OnCalendar=*-*-* 20:02
OnCalendar=*-*-* 20:04
OnCalendar=*-*-* 20:06
OnCalendar=*-*-* 20:08
OnCalendar=*-*-* 20:10
OnCalendar=*-*-* 20:12
OnCalendar=*-*-* 20:14
OnCalendar=*-*-* 20:16
OnCalendar=*-*-* 20:18
OnCalendar=*-*-* 20:20
OnCalendar=*-*-* 20:22
OnCalendar=*-*-* 20:24
OnCalendar=*-*-* 20:26
OnCalendar=*-*-* 20:28
OnCalendar=*-*-* 20:30

# 中午指定拍摄时间
OnCalendar=*-*-* 12:58
OnCalendar=*-*-* 12:59
OnCalendar=*-*-* 13:00
OnCalendar=*-*-* 13:01
OnCalendar=*-*-* 13:02
OnCalendar=*-*-* 13:03

Persistent=true
Unit=timelapse.service

[Install]
WantedBy=timers.target
